# scripts/ai.openclaw.gateway-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway-control</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/fan/.openclaw/workspace/openclaw-manager/scripts/gateway_control_worker.sh</string>
      <string>--drain</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/openclaw-manager-gateway-control.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw-manager-gateway-control.err.log</string>
  </dict>
</plist>
